# debian/DEBIAN/prerm
#!/bin/bash
systemctl stop amd-metrics-exporter.service
systemctl stop gpuagent.service

systemctl disable gpuagent.service
systemctl disable amd-metrics-exporter.service
